# Scripts/NetCore/Installation-scripts/setup-utilities/ncached.service
[Unit]
Description=NCache Daemon
DefaultDependencies=no

[Service]
User=ncache
Group=ncache
Type=forking
ExecStart=<DOTNETHOME>/dotnet <INSTALLDIR>/bin/service/Alachisoft.NCache.Daemon.dll <DOTNETHOME> <INSTALLDIR>
ExecStop=<DOTNETHOME>/dotnet <INSTALLDIR>/bin/service/Alachisoft.NCache.Daemon.dll stop
Restart=always

[Install]
WantedBy=default.target
